# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
loglevel=info
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:tts-sidecar]
command=/app/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8765 --log-level info
directory=/app/tts-sidecar
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/app/logs/tts-sidecar.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/tts-sidecar-err.log
stderr_logfile_maxbytes=2MB
environment=MODEL_DIR="/app/models",TTS_MODEL_ID="%(ENV_TTS_MODEL_ID)s",TTS_VOICE="%(ENV_TTS_VOICE)s",TTS_SAMPLE_RATE="%(ENV_TTS_SAMPLE_RATE)s",TORCH_THREADS="2"
priority=10

[program:tolk-tts-bot]
command=dotnet /app/dotnet/TolkTtsBot.dll
directory=/app/dotnet
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/app/logs/tolk-bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/tolk-bot-err.log
stderr_logfile_maxbytes=2MB
priority=20
; Передаём PORT от Render в .NET приложение
environment=Tts__SidecarUrl="http://localhost:8765",ASPNETCORE_ENVIRONMENT="%(ENV_ASPNETCORE_ENVIRONMENT)s",Bot__Name="%(ENV_BOT_NAME)s",PLAYWRIGHT_CHROMIUM_EXECUTABLE_PATH="/usr/bin/chromium",PORT="%(ENV_PORT)s"
